# adhoc-connect/etc/adhoc-connect.cron
# adhoc-connect @@ADHOC_RELEASE@@ – Fetch dhcpd.conf and restart dhcpd when needed
ADHOC_HOME=/cdg/adhoc
# $Id$
* * * * * root ${ADHOC_HOME}/adhoc-connect/sbin/adhoc-connect.sh