# adhoc-server/etc/adhoc-server-backup.cron
# adhoc-server @@ADHOC_RELEASE@@ – Make backup of the database periodically
# $Id: adhoc-server-backup.cron 387 2014-04-25 13:42:55Z bernerus@CHALMERS.SE $
10 * * * * srvadhoc /cdg/adhoc/adhoc-server/sbin/backup_mysql.sh
