# adhoc-connect/adhoc-connect.cron
# Fetch dhcpd.conf and restart dhcpd when needed
* * * * * root /cdg/sbin/adhoc-connect.sh