# adhoc-connect/adhoc-connect.cron
# adhoc-connect @@ADHOC_RELEASE@@ – Fetch dhcpd.conf and restart dhcpd when needed
# $Id:$
* * * * * root /cdg/sbin/adhoc-connect