# adhoc-connect/adhoc-connect.cron
# Fetch dhcpf.conf and restart dhcpd when needed
* * * * * /cdg/sbin/adhoc-connect.sh
